# deploy/yiyusite.service
[Unit]
Description=Yiyu Site Node Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/yiyuSite
ExecStart=/usr/bin/node server.js
Restart=on-failure
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
